# automation/launchd/com.pedan.ruankao.daily-cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pedan.ruankao.daily-cycle</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/pedan/Documents/ruankao/run-daily-cycle.command</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/pedan/Documents/ruankao</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/pedan/Documents/ruankao/ruankao-agent/logs/daily-cycle.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/pedan/Documents/ruankao/ruankao-agent/logs/daily-cycle.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
